# templates/refresh_facts.timer.epp
<%- |
  String $oncalendar,
| -%>
#
# Managed by Puppet, DO NOT EDIT
#
[Unit]
Description=Systemd Timer for Mcollective fact refresh

[Timer]
OnCalendar=<%= $oncalendar %>
Persistent=false
RandomizedDelaySec=30

[Install]
WantedBy=timers.target
